# deploy/supervisor-electronica.conf
[program:sistema-electronica]
command=/home/sistema-electronica/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 wsgi:app
directory=/home/sistema-electronica
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/sistema-electronica/error.log
stdout_logfile=/var/log/sistema-electronica/access.log
environment=
    FLASK_ENV="production",
    SECRET_KEY="$(cat /home/sistema-electronica/.env.prod | grep SECRET_KEY | cut -d= -f2-)",
    DATABASE_URL="mysql+pymysql://electronica:electronica_pass_2025@localhost/electronica_db"
